# factory-system/supervisor/factory.conf
[group:factory]
programs=factory-worker,factory-scheduler

[program:factory-worker]
process_name=%(program_name)s_%(process_num)02d
directory=/var/www/factory-system
command=/usr/bin/php artisan queue:work redis --queue=default,notifications --sleep=3 --tries=3 --timeout=120 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stopwaitsecs=3600
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/factory-system/storage/logs/worker.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5

[program:factory-scheduler]
process_name=%(program_name)s
directory=/var/www/factory-system
command=/usr/bin/php artisan schedule:work
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/factory-system/storage/logs/scheduler.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
